# board/NanoPi-K1-Plus/overlay/etc/rc.local
#!/bin/sh

echo 0 > /dev/led/nanopi:red:pwr

echo 1 > /dev/led/nanopi:green:status

