# board/Pinebook-Pro/overlay/etc/rc.local
#!/bin/sh

echo 0 > /dev/led/green:power

echo 1 > /dev/led/red:standby

